# services/mcmyadmin/run
#!/bin/bash
cd /minecraft
exec /sbin/setuser abc ./MCMA2_Linux_x86_64 